# orloy_app.service
[Unit]
Description=Orloy App
# Start after networking and Bluetooth are available
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=david
Group=david
WorkingDirectory=/home/david/Projects/orloy_app

# Activate the virtual environment, then run main.py
ExecStart=/home/david/Projects/orloy_app/venv/bin/python /home/david/Projects/orloy_app/main.py

# Restart automatically on crash (not on intentional shutdown)
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=orloy_app

[Install]
WantedBy=multi-user.target
